# supervisor.conf
[program:zifbin]
command = gunicorn app:app -b unix:///tmp/zifbin.sock -w 9
directory = /srv/zifb.in/
user = www-data
